# files/sync-feeds.timer
[Unit]
Description=start feed sync

[Timer]
OnActiveSec=24h
OnUnitActiveSec=24h
Unit=sync-feeds.service

[Install]
WantedBy=multi-user.target
